# docker/supervisord.conf
; Supervisor: corre el panel web Y el queue worker en el mismo contenedor.
; Sin el queue worker, conectar una cuenta de broker se queda "en validación"
; para siempre (el job ProvisionBrokerAccount nunca se procesa).

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/tmp/supervisord.pid

; --- Panel web (Laravel) -------------------------------------------------
[program:web]
command=sh -c "php -S 0.0.0.0:${PORT:-8080} -t public"
directory=/var/www/html
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- Queue worker (procesa los jobs: validación/deploy de cuentas, etc.) --
[program:queue]
command=php artisan queue:work --sleep=3 --tries=3 --max-time=3600 --timeout=180
directory=/var/www/html
autostart=true
autorestart=true
stopwaitsecs=190
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- Scheduler (dispara las tareas programadas cada minuto) ---------------
; Reemplaza al cron del sistema: corre el scheduler de Laravel en foreground.
; Necesario para binance:poll (recarga USDT), marketplace:charge-due y
; billing:charge-platform.
[program:scheduler]
command=php artisan schedule:work
directory=/var/www/html
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
